# leetcode-tracker.service
[Unit]
Description=LeetCode Daily Submission Tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# ── Change these two paths to match your server ──
User=ubuntu
WorkingDirectory=/home/ubuntu/leetcode-tracker
ExecStart=/home/ubuntu/leetcode-tracker/venv/bin/python tracker.py
# ─────────────────────────────────────────────────

Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
